# projects/pusher/debian/postrm
#!/bin/bash

update-rc.d buddycloud-pusher remove

#DEBHELPER#