# projects/pusher/debian/postrm
#!/bin/bash

update-rc.d -f buddycloud-pusher remove

#DEBHELPER#