# projects/pusher/debian/postinst
#!/bin/bash

mv /usr/share/buddycloud-pusher/configuration.properties.example /usr/share/buddycloud-pusher/configuration.properties
update-rc.d buddycloud-pusher defaults

#DEBHELPER#
