# projects/pusher/debian/postinst
#!/bin/bash

update-rc.d buddycloud-pusher defaults

#DEBHELPER#